# dockstarter.service
[Unit]
Description=DockStarter process supervise

[Service]
ExecStart=/opt/dockstarter

[Install]
WantedBy=multi-user.target
