# dockstarter.service
[Unit]
Description=DockStarter process supervise

[Service]
ExecStart=/opt/dockstarter
Restart=always

[Install]
WantedBy=multi-user.target
