# minisv.service
[Unit]
Description=minisv process supervise

[Service]
ExecStart=/opt/minisv
Restart=always

[Install]
WantedBy=multi-user.target
